# retro.service
[Unit]
Description=Retro Board
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=retro
WorkingDir=/opt/retro
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=DATABASE_URL=postgresql://retro:retro@localhost:5432/retro
Environment=ORIGIN=http://localhost
ExecStartPre=/usr/bin/node migrate.js
ExecStart=/usr/bin/node --max-old-space-size=256 server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
